# systemd/script-manager-ui.service
[Unit]
Description=Script Manager UI
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=thomas
Group=thomas
WorkingDirectory=/opt/apps/script-manager-ui
Environment="PATH=/opt/apps/script-manager-ui/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/apps/script-manager-ui/venv/bin/python main.py
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/script-manager-ui /var/log/honeyfile_monitor.log /var/lib/honeyfile_last_alert_ts

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=script-manager-ui

[Install]
WantedBy=multi-user.target
